# docker/root/etc/services.d/porn-vault/run
#!/usr/bin/with-contenv bash

echo "Starting Porn Vault"

cd /
exec s6-setuidgid abc /porn-vault